# rtl433log_supervisor.conf
[program:rtl433log]
command=python /var/www/app/rtl433log/rtl433log_threading.py
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/log/rtl433log_threding.out.log
stderr_logfile=/var/log/rtl433log_threding.err.log
